# debian/cron.d/intelmq-update-data
# /etc/cron.d/intelmq-update-data: crontab fragment for intelmq
#  This updates the data files used by some expert bots.
#
# m h  dom mon dow      command
# Update data for tor_nodes bot:
11 0 * * *     intelmq   /usr/bin/intelmq-update-tor-nodes /var/lib/intelmq/bots/tor_nodes/tor_nodes.dat
# Update data for maxmind_geoip bot:
17 0 * * *     intelmq   /usr/bin/intelmq-update-geoip-data /var/lib/intelmq/bots/maxmind_geoip/GeoLite2-City.mmdb
# Update data for asn_lookup bot:
23 0 * * *     intelmq   /usr/bin/intelmq-update-asn-data /var/lib/intelmq/bots/asn_lookup/ipasn.dat
